# meta-facebook/meta-yosemite4/recipes-yosemite4/pldm-monitor/files/pldm-monitor.service
[Unit]
Description=Monitor pldm sensors
ConditionPathExists=/run/openbmc/pldm_monitor_requested

[Service]
Type=oneshot
ExecStart=/usr/libexec/pldm-monitor

[Install]
WantedBy=multi-user.target

