# meta-facebook/meta-yosemite4/recipes-yosemite4/nic-reset/files/nic-reset@.service
[Unit]
Description=Reset specified NIC

[Service]
Type=oneshot
ExecStart=/usr/libexec/nic-reset %i
