# meta-facebook/meta-yosemite4/recipes-yosemite4/fan-manual-mode/files/fan-manual-mode.service
[Unit]
Description=Set fans to manual mode if environment variable exists

[Service]
Type=oneshot
ExecStart=/usr/libexec/fan-manual-mode

[Install]
WantedBy=multi-user.target
