# meta-facebook/meta-yosemite4/recipes-yosemite4/managed-reboot/files/managed-reboot.service
[Unit]
Description=Request CHAI for a reboot via SEL

[Service]
Type=oneshot
ExecStart=/usr/libexec/managed-reboot
