# meta-facebook/meta-yosemite4/recipes-phosphor/health/files/mem-dump.service
[Unit]
Description=Dump stats on who is using up memory

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-health-monitor/mem-dump
SyslogIdentifier=mem-dump
