# common/recipes-core/rackmon2/rackmon/rackmond.service.in
# Copyright 2014-present Facebook. All Rights Reserved.
[Unit]
Description=Start rackmond
REPLACE_WITH_SERVICE_AFTER
REPLACE_WITH_SERVICE_WANT

[Service]
ExecStart=/usr/local/bin/rackmond REPLACE_WITH_INTERFACE_CONFIG_PATH REPLACE_WITH_REGMAP_CONFIG_DIR_PATH

[Install]
WantedBy=multi-user.target
